# gnome-2-24/daemon/cdda.mount.in
[Mount]
Type=cdda
Exec=@libexecdir@/gvfsd-cdda
AutoMount=false
