# gnome-2-24/daemon/computer.mount.in
[Mount]
Type=computer
Exec=@libexecdir@/gvfsd-computer
AutoMount=true

