# gnome-2-24/daemon/dns-sd.mount.in
[Mount]
Type=dns-sd
Exec=@libexecdir@/gvfsd-dnssd
AutoMount=true
