# gnome-2-24/daemon/gphoto2.mount.in
[Mount]
Type=gphoto2
Exec=@libexecdir@/gvfsd-gphoto2
AutoMount=false
