# gnome-2-24/daemon/localtest.mount.in
[Mount]
Type=localtest
Exec=@libexecdir@/gvfsd-localtest
AutoMount=false
