# gnome-2-24/daemon/obexftp.mount.in
[Mount]
Type=obex
Exec=@libexecdir@/gvfsd-obexftp
AutoMount=false
