# gnome-2-24/daemon/smb.mount.in
[Mount]
Type=smb-share
Exec=@libexecdir@/gvfsd-smb
AutoMount=false
Scheme=smb
