# metadata/gvfs-metadata.service.in
[Unit]
Description=Virtual filesystem metadata service
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfsd-metadata
Type=dbus
BusName=org.gtk.vfs.Metadata
Slice=session.slice
