# monitor/afc/gvfs-afc-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - Apple File Conduit monitor
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfs-afc-volume-monitor
Type=dbus
BusName=org.gtk.vfs.AfcVolumeMonitor
Slice=session.slice
