# daemon/afc.mount.in
[Mount]
Type=afc
Exec=@libexecdir@/gvfsd-afc
AutoMount=false
Scheme=afc
DefaultPort=1

