# daemon/afp-browse.mount.in
[Mount]
Type=afp-server
Exec=@libexecdir@/gvfsd-afp-browse
AutoMount=false
Scheme=afp
DefaultPort=548
HostnameIsInetAddress=true

