# daemon/afp.mount.in
[Mount]
Type=afp-volume
Exec=@libexecdir@/gvfsd-afp
AutoMount=false
Scheme=afp
DefaultPort=548
HostnameIsInetAddress=true