# daemon/archive.mount.in
[Mount]
Type=archive
Exec=@libexecdir@/gvfsd-archive
AutoMount=false
Scheme=archive
