# daemon/burn.mount.in
[Mount]
Type=burn
Exec=@libexecdir@/gvfsd-burn
AutoMount=true

