# daemon/dav+sd.mount.in
[Mount]
Type=dav+sd;davs+sd
Exec=@libexecdir@/gvfsd-dav
AutoMount=false
