# daemon/dav.mount.in
[Mount]
Type=dav
Exec=@libexecdir@/gvfsd-dav
AutoMount=false
