# daemon/dns-sd.mount.in
[Mount]
Type=dns-sd
Exec=@libexecdir@/gvfsd-dnssd
AutoMount=true
DBusName=org.gtk.vfs.mountpoint_dnssd
