# daemon/ftp.mount.in
[Mount]
Type=ftp
Exec=@libexecdir@/gvfs-daemon-ftp
AutoMount=false
