# daemon/ftp.mount.in
[Mount]
Type=ftp
Exec=@libexecdir@/gvfsd-ftp
AutoMount=false
Scheme=ftp
DefaultPort=21
HostnameIsInetAddress=true
