# daemon/ftp.mount.in
[Mount]
Type=ftp
Exec=@libexecdir@/gvfsd-ftp
AutoMount=false
