# daemon/ftpes.mount.in
[Mount]
Type=ftpes
Exec=@libexecdir@/gvfsd-ftp
AutoMount=false
Scheme=ftpes
DefaultPort=21
HostnameIsInetAddress=true
