# daemon/ftpis.mount.in
[Mount]
Type=ftpis
Exec=@libexecdir@/gvfsd-ftp
AutoMount=false
Scheme=ftpis
DefaultPort=990
HostnameIsInetAddress=true
