# daemon/ftps.mount.in
[Mount]
Type=ftps
Exec=@libexecdir@/gvfsd-ftp
AutoMount=false
Scheme=ftps
DefaultPort=21
HostnameIsInetAddress=true
