# daemon/google.mount.in
[Mount]
Type=google-drive
Exec=@libexecdir@/gvfsd-google
AutoMount=false
