# daemon/gvfs-daemon.service.in
[Unit]
Description=Virtual filesystem service

[Service]
ExecStart=@libexecdir@/gvfsd
Type=dbus
BusName=org.gtk.vfs.Daemon
