# daemon/gvfs-daemon.service.in
[Unit]
Description=Virtual filesystem service
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfsd
Type=dbus
BusName=org.gtk.vfs.Daemon
Slice=session.slice
