# daemon/http.mount.in
[Mount]
Type=http
Exec=@libexecdir@/gvfsd-http
AutoMount=true
DBusName=org.gtk.vfs.mountpoint.http
