# daemon/mtp.mount.in
[Mount]
Type=mtp
Exec=@libexecdir@/gvfsd-mtp
AutoMount=false
