# daemon/nfs.mount.in
[Mount]
Type=nfs
Exec=@libexecdir@/gvfsd-nfs
AutoMount=false
Scheme=nfs
