# daemon/onedrive.mount.in
[Mount]
Type=onedrive
Exec=@libexecdir@/gvfsd-onedrive
AutoMount=false
