# daemon/recent.mount.in
[Mount]
Type=recent
Exec=@libexecdir@/gvfsd-recent
AutoMount=true

