# daemon/sftp.mount.in
[Mount]
Type=sftp
Exec=@libexecdir@/gvfsd-sftp
AutoMount=false
