# daemon/sftp.mount.in
[Mount]
Type=ftp
Exec=@libexecdir@/gvfs-daemon-sftp
AutoMount=false
