# daemon/sftp.mount.in
[Mount]
Type=sftp
Exec=@libexecdir@/gvfs-daemon-sftp
AutoMount=false
