# daemon/smb-browse.mount.in
[Mount]
Type=smb-network;smb-server
Exec=@libexecdir@/gvfsd-smb-browse
DBusName=org.gtk.vfs.mountpoint.smb_browse
AutoMount=false
Scheme=smb
