# daemon/smb-browse.mount.in
[Mount]
Type=smb-network;smb-server
Exec=@libexecdir@/gvfs-daemon-smb-browse --mount
DBusName=org.gtk.vfs.mountpoint.smb_browse
ObjPath=/org/gtk/vfs/mountpoint/smb_browse
AutoMount=true

