# daemon/smb.mount.in
[Mount]
Type=smb-share
Exec=@libexecdir@/gvfs-daemon-smb
AutoMount=false
