# daemon/trash.mount.in
[Mount]
Type=trash
Exec=@libexecdir@/gvfs-daemon-trash
AutoMount=true

