# daemon/wsdd.mount.in
[Mount]
Type=wsdd
Exec=@libexecdir@/gvfsd-wsdd
AutoMount=true