# monitor/gphoto2/gvfs-gphoto2-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - digital camera monitor
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfs-gphoto2-volume-monitor
Type=dbus
BusName=org.gtk.vfs.GPhoto2VolumeMonitor
Slice=session.slice
