# monitor/hal/gvfs-hal-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - deprecated HAL monitor

[Service]
ExecStart=@libexecdir@/gvfs-hal-volume-monitor
Type=dbus
BusName=org.gtk.vfs.HalVolumeMonitor
