# monitor/udisks2/gvfs-udisks2-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - disk device monitor

[Service]
ExecStart=@libexecdir@/gvfs-udisks2-volume-monitor
Type=dbus
BusName=org.gtk.vfs.UDisks2VolumeMonitor
