# monitor/mtp/gvfs-mtp-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - Media Transfer Protocol monitor
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfs-mtp-volume-monitor
Type=dbus
BusName=org.gtk.vfs.MTPVolumeMonitor
Slice=session.slice
