# monitor/mtp/gvfs-mtp-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - Media Transfer Protocol monitor

[Service]
ExecStart=@libexecdir@/gvfs-mtp-volume-monitor
Type=dbus
BusName=org.gtk.vfs.MTPVolumeMonitor
