# trunk/daemon/dav.mount.in
[Mount]
Type=dav;davs
Exec=@libexecdir@/gvfsd-dav
AutoMount=false
