# trunk/daemon/network.mount.in
[Mount]
Type=network
Exec=@libexecdir@/gvfsd-network
AutoMount=true
