# trunk/daemon/sftp.mount.in
[Mount]
Type=sftp
Exec=@libexecdir@/gvfsd-sftp
AutoMount=false
Scheme=sftp
SchemeAliases=ssh
DefaultPort=22
HostnameIsInetAddress=true
