# trunk/daemon/trash.mount.in
[Mount]
Type=trash
Exec=@libexecdir@/gvfsd-trash
AutoMount=true

