# monitor/goa/gvfs-goa-volume-monitor.service.in
[Unit]
Description=Virtual filesystem service - GNOME Online Accounts monitor
PartOf=graphical-session.target

[Service]
ExecStart=@libexecdir@/gvfs-goa-volume-monitor
Type=dbus
BusName=org.gtk.vfs.GoaVolumeMonitor
Slice=session.slice
